# repo.crontab
01 00 * * *   root /usr/lib/obs/service/repo-service-cronjob
